# ADC_Volume_Control/adc-volume.service
[Unit]
Description=ADC volume control

[Service]
User=pi
Group=pi
ExecStart=/home/pi/bin/adc-volume-control

[Install]
WantedBy=multi-user.target